# files/node-exporter-pacman.service
[Unit]
Description=Run pacman metrics checkin

[Service]
Type=simple
ExecStart=/usr/local/bin/node-exporter-pacman.sh
