# files/node-exporter@.timer
[Unit]
Description=Schedule %i metric checkins

[Timer]
OnBootSec=5min
OnUnitActiveSec=600

[Install]
WantedBy=multi-user.target
